# noctalia/.config/noctalia/plugins/battery-threshold/99-battery-threshold.rules
# Battery Threshold Control - udev rule
# Grants write access to charge_control_end_threshold for users in the
# 'battery_ctl' group, only for BAT0
SUBSYSTEM=="power_supply", KERNEL=="BAT0", \
    RUN+="/bin/chgrp battery_ctl /sys$devpath/charge_control_end_threshold", \
    RUN+="/bin/chmod g+w /sys$devpath/charge_control_end_threshold"

